# supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:bench-web]
command=/home/frappe/.local/bin/bench serve --port 8000 --host 0.0.0.0
directory=/home/frappe/frappe-bench
user=frappe
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/bench-web.log
stderr_logfile=/var/log/supervisor/bench-web.log
priority=100

[program:bench-schedule]
command=/home/frappe/.local/bin/bench schedule
directory=/home/frappe/frappe-bench
user=frappe
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/bench-schedule.log
stderr_logfile=/var/log/supervisor/bench-schedule.log
priority=200

[program:bench-worker-default]
command=/home/frappe/.local/bin/bench worker --queue default
directory=/home/frappe/frappe-bench
user=frappe
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/bench-worker-default.log
stderr_logfile=/var/log/supervisor/bench-worker-default.log
priority=300

[program:bench-worker-long]
command=/home/frappe/.local/bin/bench worker --queue long
directory=/home/frappe/frappe-bench
user=frappe
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/bench-worker-long.log
stderr_logfile=/var/log/supervisor/bench-worker-long.log
priority=400

[program:bench-worker-short]
command=/home/frappe/.local/bin/bench worker --queue short
directory=/home/frappe/frappe-bench
user=frappe
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/bench-worker-short.log
stderr_logfile=/var/log/supervisor/bench-worker-short.log
priority=500

[program:socketio]
command=/home/frappe/.local/bin/bench node-socketio
directory=/home/frappe/frappe-bench
user=frappe
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/socketio.log
stderr_logfile=/var/log/supervisor/socketio.log
priority=600